# data/init/systemd/default/process.service.erb
[Unit]
StopWhenUnneeded=true
Requires=<%= name %>-<%= process_name %>.service
After=<%= name %>-<%= process_name %>.service

[Service]
Environment=PORT=PORT_NUM
ExecStart=/usr/bin/<%= name %> run <%= process_name %> >> /var/log/<%= name %>/<%= rocess_name %>-PROCESS_NUM.log 2>&1
Restart=always

[Install]
WantedBy=<%= name %>-<%= process_name %>.service
